# systemd/tryx-panorama.service
[Unit]
Description=TRYX Panorama SE 360 Display Manager
Documentation=https://github.com/fadli0029/reed-tpse
After=local-fs.target

[Service]
Type=simple
ExecStart=%h/.local/bin/tryx-panorama-manager
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
